# sms-bridge.timer
[Unit]
Description=Bridge SMSes to matrix

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
